# chromebook-audio.service
[Unit]
Description=Symlinks config files for the current Chromebook on boot.
Before=systemd-sysext.service

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=chromebook-audio

[Install]
WantedBy=default.target
